# site_tracker/officetime.conf
start on startup

chdir /home/vb/projects/git_hosted/site_tracker
exec phantomjs index.js
